# udev/rules.d/99-notify-lowbat.rules
SUBSYSTEM=="power_supply", ATTR{status}=="Discharging", ATTR{capacity}=="15", RUN+="/etc/udev/notify-lowbat.sh"
